# tailscale-cert.timer.j2
[Unit]
Description=Renew Tailscale cert

[Timer]
OnCalendar=monthly
Unit=%i.service
Persistent=true

[Install]
WantedBy=timers.target
